# dev/3ginetgw.crontab
*/3 * * * *	root	/sbin/3ginetgw.sh chkconn >> /etc/3ginetgw/3ginetgw.log
